# docker/rootfs/etc/s6-overlay/s6-rc.d/svc-redis/finish
#!/command/with-contenv /bin/bash
# shellcheck shell=bash

pkill -f 'redis-server'