# docker/rootfs/etc/s6-overlay/s6-rc.d/init-015-bootstrap-login/up
#!/command/execlineb -P
foreground { mkdir -p /run/s6-rc }
s6-touch /run/s6-rc/init-015-bootstrap-login
