# docker/rootfs/etc/s6-overlay/s6-rc.d/init-010-update/run
#!/command/with-contenv /bin/bash
# shellcheck shell=bash

info() {
    echo -e "[\033[32mINFO\033[0m] $1"
}

warn() {
    echo -e "[\033[33mWARN\033[0m] $1"
}

if [ "${NASTOOL_IMMUTABLE_IMAGE:-false}" = "true" ]; then
    info "当前镜像采用不可变部署，跳过容器内自更新；如需升级，请重新构建并部署新镜像。"
    if [ "${NASTOOL_AUTO_UPDATE:-false}" = "true" ]; then
        warn "已忽略 NASTOOL_AUTO_UPDATE=true；不可变镜像不支持容器内 git 更新。"
    fi
elif [ "${NASTOOL_AUTO_UPDATE:-false}" = "true" ]; then
    warn "当前运行环境未启用不可变镜像，但容器内自更新链路已不再推荐使用；建议改为重新构建并部署新镜像。"
else
    info "程序自动升级已关闭。"
fi
